# systemd/agents-reconcile.timer
[Unit]
Description=Periodically reconcile ~/.config/agents with origin/main

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
