# systemd/otel-stack.service
[Unit]
Description=agents local OpenTelemetry stack (collector + Prometheus + Grafana)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:%h/.local/bin:/usr/bin:/bin
WorkingDirectory=%h/.config/agents/observability
ExecStart=/usr/bin/env docker compose -f %h/.config/agents/observability/docker-compose.yml up -d
ExecStop=/usr/bin/env docker compose -f %h/.config/agents/observability/docker-compose.yml down

[Install]
WantedBy=default.target
